# twapake.cron
#!/bin/sh

#
# print an every day twapake snapshot in /var/log/twapake/
#

/usr/sbin/twapake --logfile
